# transports/mqtt/hyperdrive-transport-mqtt.service.in
[Unit]
Description=Hyperdrive MQTT Transport
Requires=hyperdrive.service
After=hyperdrive.service
After=network.target

[Service]
ExecStart=@INSTALL_BIN_DIR@/hyperdrive-transport-mqtt
User=mqtt-transport
