# transports/http/hyperdrive-transport-http.service.in
[Unit]
Description=Hyperdrive HTTP Transport
Requires=hyperdrive.service
After=hyperdrive.service
After=sockets.target

[Service]
ExecStart=@INSTALL_BIN_DIR@/hyperdrive-transport-http
User=http-transport
