# discoveryservices/hyperdiscovery/hyperdrive-discoveryservice-hyperdiscovery.socket
[Unit]
Description=Hyperdrive Hyperdiscovery Discovery service socket

[Socket]
ListenDatagram=44389

[Install]
WantedBy=sockets.target
